# config/autostart/gnome-terminnal.desktop
[Desktop Entry]
Exec=gnome-terminal --geometry=94x47+0+0
Type=Application
